# disk/usbnosleep/60-persistent-usbstorage-sleep.rules
ACTION=="add|change", KERNEL=="sd*[!0-9]", SUBSYSTEMS=="usb", RUN+="/usr/bin/sdparm --clear STANDBY /dev/%k"
